# systemd/user/onedrive.service.d/override.conf
[Service]
Restart=on-failure
RestartSec=3
RestartPreventExitStatus=3
CPUWeight=20
StartupCPUWeight=40
CPUShares=20
CPUQuota=20%
IOWeight=20
MemorySwapMax=0
